# config/django-server.conf
[program:django-server]
command=/bin/bash -c "exec gunicorn panelProject.wsgi:application --bind 0.0.0.0:8000 --workers ${GUNICORN_WORKERS:-2} --timeout ${GUNICORN_TIMEOUT:-120} --access-logfile /home/python/panel/logs/gunicorn-access.log --error-logfile /home/python/panel/logs/gunicorn-error.log"
directory=/home/python/panel
user=root
autostart=true
autorestart=true
stdout_logfile=/home/python/panel/logs/django.log
stderr_logfile=/home/python/panel/logs/django.err
environment=PYTHONUNBUFFERED=1
